# docker/devbox/common/supervisord-phpfpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm7.4 -F
autostart=true
autorestart=unexpected
stdout_logfile=/var/log/phpfpm/stdout.log
stderr_logfile=/var/log/phpfpm/stderr.log