# docker/ubuntu_precice_nginx/common/supervisord-phpfpm.conf
[program:phpfpm]
command=/start-phpfpm.sh
numprocs=1
autostart=true
autorestart=true